# com.sour4bh.cc-controller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sour4bh.cc-controller</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/sour4bh/resume/cc-controller/venv/bin/python</string>
        <string>-m</string>
        <string>cc_controller.daemon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/sour4bh/resume/cc-controller</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/sour4bh/resume/cc-controller</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/sour4bh/Library/Logs/cc-controller.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/sour4bh/Library/Logs/cc-controller.error.log</string>

    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
